# deploy/supervisor/backend-scheduler.conf
[program:backend-scheduler]
command=/home/ubuntu/fastapi-async-sqlmodel-boilerplate/deploy/supervisor/backend-scheduler
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true
stderr_logfile=/home/ubuntu/fastapi-async-sqlmodel-boilerplate/deploy/logs/scheduler-error.log
stdout_logfile=/home/ubuntu/fastapi-async-sqlmodel-boilerplate/deploy/logs/scheduler-output.log